# com.example.invoice-attachment-archive.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.invoice-attachment-archive</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/path/to/notion-slack-automation/jobs/invoice-attachment-archive/index.mjs</string>
    <string>--apply</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/path/to/notion-slack-automation</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/path/to/notion-slack-automation/logs/invoice-attachment-archive.out.log</string>

  <key>StandardErrorPath</key>
  <string>/path/to/notion-slack-automation/logs/invoice-attachment-archive.err.log</string>
</dict>
</plist>
